# .config/autostart/redshift.desktop
[Desktop Entry]
Exec=redshift -l 50:1
Name=redshift
Type=Application
Version=1.0
